# systemd/akshare-mirror.timer
[Unit]
Description=Trigger AKShare -> GitHub mirror nightly

[Timer]
# 01:30 CST — after mainland market close, before slate nightly at 02:00 US.
# Adjust via `sudo systemctl edit akshare-mirror.timer` if needed.
OnCalendar=*-*-* 01:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
